# deploy/systemd/ai-trader.service
[Unit]
Description=AI Trader Bot + Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=__APP_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=/bin/bash -lc 'source .venv/bin/activate && python -m ai_trader_bot --log-level INFO'
Restart=always
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
